# core/maintenance/reboot-complete-m8.service
[Unit]
Description=Publish boot-complete notice to JCTsh log dashboard (M8 photo-server)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/jctsh/heartbeat.env
ExecStart=/usr/bin/mosquitto_pub -h 192.168.1.117 -p 1883 -u photo-server -P ${MQTT_PASSWORD} -t jctsh/server/photo-server/log -m '{"component":"photo-server","category":"System","message":"Boot complete."}'

[Install]
WantedBy=multi-user.target
