# core/maintenance/scheduled-reboot.service
[Unit]
Description=Weekly scheduled maintenance reboot

[Service]
Type=oneshot
ExecStart=/sbin/reboot
